# Files from DROP/etc/rc.local
wifi 
hostapd_cli -p /var/run/hostapd-phy0 karma_disable
#Don't touch anything above this line


#Add your commands above this
exit 0
