# Files from DROP/etc/hotplug.d/button/buttons
#!/bin/sh
logger $BUTTON
logger $ACTION
